# config/check-smart.service
[Unit]
Description=Daily SMART health check

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/check-smart.sh
