# config/sync-root-b.timer
[Unit]
Description=Weekly root_B sync timer

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
