# config/sync-root-c.timer
[Unit]
Description=Monthly root_C sync timer

[Timer]
OnCalendar=*-*-01 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
